# systemd/keyserver-watcher.path
[Path]
PathModified=/opt/keyserver/keys/mrdl_ubuntu.pub
PathModified=/opt/keyserver/keys/mrdl_ca.pub

[Install]
WantedBy=multi-user.target