# systemd/kscreate-watcher.path
[Unit]
Description=keyserver path

[Path]
Unit=kscreate.service
PathModified=/home/dnlennon/.ssh/keys/mrdl/mrdl_ubuntu.pub
PathModified=/home/dnlennon/.ssh/keys/mrdl/mrdl_ca.pub

[Install]
WantedBy=kscreate.service